# systemd/kiosk-restart.service
[Unit]
Description=Restart kiosk browser and backend
Documentation=https://github.com/yourusername/familywall
After=network.target

[Service]
Type=oneshot
User=jassmith
ExecStart=/home/jassmith/familywall/scripts/restart-kiosk.sh

[Install]
WantedBy=multi-user.target
